# assets/cleanup-pastes.timer
[Unit]
Description=Cleanup old pastes

[Timer]
OnUnitActiveSec=3h

[Install]
WantedBy=timers.target
